# supervisord.conf
[supervisord]
nodaemon=true

[program:libvirtd]
command=/usr/sbin/libvirtd

[program:virtlockd]
command=/usr/sbin/virtlockd

[program:virtlogd]
command=/usr/sbin/virtlogd